# FV3GFSwfm/test_GSDnoah_l128/test_GSDnoah_l128.crontab

#################### test_GSDnoah_l128 ####################
MAILTO=""
*/5 * * * * /apps/rocoto/1.3.2/bin/rocotorun -d /scratch2/BMC/gsd-fv3-dev/Judy.K.Henderson/test/gsd_ccpp_v16b_restructure_update/FV3GFSwfm/test_GSDnoah_l128/test_GSDnoah_l128.db -w /scratch2/BMC/gsd-fv3-dev/Judy.K.Henderson/test/gsd_ccpp_v16b_restructure_update/FV3GFSwfm/test_GSDnoah_l128/test_GSDnoah_l128.xml
#################################################################

